# templates/default/runsvdir.conf.erb
# generated by limits cookbook
# if you get Chef::Mixin::Template::TemplateError here, stop using ruby 1.8 and use omnibus
<% @filtered_limits.each_key do |type| %>
  <% hard = @filtered_limits[type]['hard'] %>
  <% soft = @filtered_limits[type]['soft'] %>
limit <%= type %> <%= soft %> <%= hard %>
<% end %>
kill timeout 5
pre-stop script
  set +e
  killall -HUP runsvdir
  exit 0
end script
pre-start script
  set +e
  killall -TERM runsv
  sleep 3
  killall -KILL runsv
  exit 0
end script
normal exit 111 0 TERM
start on runlevel [2345]
stop on shutdown
respawn
exec /usr/sbin/runsvdir-start
